# cloudflareargo/rootfs/etc/services.d/argo/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Community Hass.io Add-ons: Node-RED
# Starts Node-RED
# ==============================================================================
declare -a options
declare certificate 
declare url

url=$(bashio::config 'url')
certificate=$(bashio::config 'certificate')

bashio::log.info "Starting CloudFlare Argo"

#Check URL
    attempt_counter=0
    max_attempts=150
until $(curl --output /dev/null --silent --head --fail ${url}); do
    if [ ${attempt_counter} -eq ${max_attempts} ];then
      echo "Max attempts reached"
      exit 1
    fi

    printf '.'
    attempt_counter=$(($attempt_counter+1))
    sleep 5
done
# / Check URL
exec /opt/cloudflared --config /etc/cloudflared/1.yml --origincert ${certificate}
