# addon-cloudflare-argo/cloudflareargo/rootfs/etc/services.d/argo2/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Community Hass.io Add-ons: Node-RED
# Starts Node-RED
# ==============================================================================
declare -a options
declare certificate2
declare url2


url2=$(bashio::config 'url2')
certificate2=$(bashio::config 'certificate2')

if bashio::config.has_value 'hostname2'; then
bashio::log.info "Starting CloudFlare Argo Tunnel 2"

# / Check URL
exec /opt/cloudflared --config /etc/cloudflared/2.yml --origincert ${certificate2} --no-tls-verify --proxy-connect-timeout 150s
fi

if [[ $(bashio::config "hostname2") == "null" ]]; then
SUPERVISED_DIR=/var/run/s6/services

s6-svc -x ${SUPERVISED_DIR}/argo2
rm -rf ${SUPERVISED_DIR}/argo2
s6-svscanctl -a ${SUPERVISED_DIR}       
fi
