# addon-cloudflare-argo/cloudflareargo/rootfs/etc/services.d/argo3/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Community Hass.io Add-ons: Node-RED
# Starts Node-RED
# ==============================================================================
declare -a options
declare certificate3
declare url3


url3=$(bashio::config 'url3')
certificate3=$(bashio::config 'certificate3')


if bashio::config.has_value 'hostname3'; then
bashio::log.info "Starting CloudFlare Argo Tunnel 3"

# / Check URL
exec /opt/cloudflared --config /etc/cloudflared/3.yml --origincert ${certificate3} --no-tls-verify --proxy-connect-timeout 150s
fi

if [[ $(bashio::config "hostname3") == "null" ]]; then
SUPERVISED_DIR=/var/run/s6/services

s6-svc -x ${SUPERVISED_DIR}/argo3
rm -rf ${SUPERVISED_DIR}/argo3
s6-svscanctl -a ${SUPERVISED_DIR}       
fi
