# addon-tuyagrab/tuyagrab/rootfs/etc/services.d/tuyagrab/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Community Hass.io Add-ons: Node-RED
# Starts Node-RED
# ==============================================================================
declare -a options
declare tuyakey
declare tuyasecret
declare tuyadevice
declare tuyaregion
declare refreshtime

tuyakey=$(bashio::config 'tuyakey')
tuyasecret=$(bashio::config 'tuyasecret')
tuyadevice=$(bashio::config 'tuyadevice')
tuyaregion=$(bashio::config 'tuyaregion')
refreshtime=$(bashio::config 'refreshtime')

bashio::log.info "Starting Tuya Device Capture"

json='{"apiKey": "'"${tuyakey}"'", "apiDeviceID": "'"${tuyadevice}"'", "apiSecret": "'"${tuyasecret}"'", "apiRegion": "'"${tuyaregion}"'" }'
cd /opt
echo $json > tinytuya.json
(echo "y"; sleep 10; echo "n"; sleep 30; echo "n") |  python3 -m tinytuya wizard 2>&1
cp /opt/tuya-raw.json /config/tuyadevices.txt 

sleep ${refreshtime}
