# addon-tuyagrab/tuyagrab/rootfs/etc/services.d/tuyagrab/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Community Hass.io Add-ons: Node-RED
# Starts Node-RED
# ==============================================================================
declare -a options
declare tuyakey
declare tuyasecret
declare tuyadevice

tuyakey=$(bashio::config 'tuyakey')
tuyasecret=$(bashio::config 'tuyasecret')
tuyadevice=$(bashio::config 'tuyadevice')

bashio::log.info "Starting Tuya Device Capture"

#exec ( 
(echo 'n'; sleep 5; echo 'xufadfrua4ji0pelei4k'; sleep 5; echo 'f6583ac0332a461694d712e067b03ea0'; sleep 5; echo '40374727c82b96ec9e1c') | exec tuya-cli wizard

sleep 60m
