# addon-tuyagrab/tuyagrab/rootfs/etc/services.d/tuyagrab/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Community Hass.io Add-ons: Node-RED
# Starts Node-RED
# ==============================================================================
declare -a options
declare tuyakey
declare tuyasecret
declare tuyadevice

tuyakey=$(bashio::config 'tuyakey')
tuyasecret=$(bashio::config 'tuyasecret')
tuyadevice=$(bashio::config 'tuyadevice')

bashio::log.info "Starting Tuya Device Capture"

json='{ "apiKey": "xufadfrua4ji0pelei4k", "apiDeviceID": "40374727c82b96ec9e1c", "apiSecret": "f6583ac0332a461694d712e067b03ea0", "apiRegion": "eu" }'
cd /opt
echo $json > tinytuya.json
(echo "y"; sleep 10; echo "n") |  python3 -m tinytuya wizard
