# cloudflareargo/rootfs/etc/services.d/argo3/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Community Hass.io Add-ons: Node-RED
# Starts Node-RED
# ==============================================================================
declare -a options
sleep 40


if bashio::config.has_value 'hostname3'; then
    bashio::log.info "Starting CloudFlare Argo Tunnel 3"
    exec /opt/cloudflared --config /etc/cloudflared/3.yml --no-autoupdate
fi
