# addon-tuyagrab/tuyagrab/rootfs/etc/services.d/autoshutdown/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Community Hass.io Add-ons: Node-RED
# Starts Node-RED
# ==============================================================================
declare -a options
declare shutdowntime

shutdowntime=$(bashio::config 'shutdowntime')

bashio::log.info "Starting Auto Shutdown Timer"

time=$(( ${shutdowntime} + 2 ))
sleep ${time}

shutdown -rf now
