# addon-cloudflare-argo/cloudflareargo/rootfs/etc/services.d/argo/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Community Hass.io Add-ons: Cloudflare-argo
# Starts CloudFlare Argo
# ==============================================================================
declare -a options
declare certificate 

certificate=$(bashio::config 'certificate')


bashio::log.info "Starting CloudFlare Argo"

# / Check URL
exec /opt/cloudflared --no-tls-verify --origincert=${certificate} tunnel --config=/config/cf-argo/config.yml run

Sleep 15
